# root/etc/services.d/nginx/run
#!/usr/bin/with-contenv sh

exec nginx -c /config/nginx.conf -g "daemon off;"
